# recipes/mosquitto/files/mosquitto.service
[Unit]
Description=Mosquitto - MQTT broker

[Service]
ExecStart=/usr/sbin/mosquitto
Restart=on-failure

[Install]
WantedBy=multi-user.target
